# IPv6_disabled/files/etc/uci-defaults/999-QOL_config
uci -q batch << EOI

set network.globals.packet_steering='2'

# Wireless radio0 specific settings
set wireless.radio0.htmode='HE40'
set wireless.radio0.channel='auto'
set wireless.radio0.country='PT'

# Wireless radio1 specific settings
set wireless.radio1.htmode='HE160'
set wireless.radio1.channel='auto'
set wireless.radio1.country='PT'

# Firewall settings
set firewall.@zone[1].input='DROP'
set firewall.@zone[1].forward='DROP'

# uHTTPd settings
set uhttpd.main.redirect_https='1'

# Disable IPv6
set network.lan.ipv6='0'
set network.wan.ipv6='0'
delete network.wan6
delete network.lan.ip6assign

# Disable RA and DHCPv6 so no IPv6 IPs are handed out
delete dhcp.lan.dhcpv6
delete dhcp.lan.ra

# Disable the LAN delegation
set network.lan.delegate='0'

# Delete the IPv6 ULA Prefix
delete network.globals.ula_prefix

EOI
